# package/wwan/app/luci-app-checknet/root/etc/uci-defaults/luci-checknet
#!/bin/sh

uci -q batch <<-EOF >/dev/null
	delete ucitrack.@checknet[-1]
	add ucitrack checknet
	set ucitrack.@checknet[-1].init=checknet
	commit ucitrack
EOF

rm -f /tmp/luci-indexcache
exit 0
